# templates/etc/systemd/system/airflow-worker.service.j2
# {{ ansible_managed }}
# Please do not modify this file manually.

[Unit]
Description=Airflow celery worker daemon
After=network.target

[Service]
Environment="AIRFLOW_HOME={{ airflow_install_directory }}"
User={{ airflow_user }}
Group={{ airflow_group }}
Type=simple
ExecStart=/bin/bash -c 'exec /usr/local/bin/airflow worker {{ airflow_worker_arguments }} 2>&1 >> {{ default_airflow.core.base_log_folder }}/airflow-worker.log'
Restart=always
RestartSec=11s

[Install]
WantedBy=multi-user.target
