# templates/etc/systemd/system/airflow-flower.service.j2
# {{ ansible_managed }}
# Please do not modify this file manually.

[Unit]
Description=Airflow celery flower
After=network.target

[Service]
Environment="AIRFLOW_HOME={{ airflow_install_directory }}"
User={{ airflow_user }}
Group={{ airflow_user }}
Type=simple
ExecStart=/bin/bash -c 'exec /usr/local/bin/airflow flower 2>&1 >> {{ default_airflow.core.base_log_folder }}/airflow-flower.log'
Restart=always
RestartSec=11s

[Install]
WantedBy=multi-user.target
