# templates/etc/systemd/system/airflow-webserver.service.j2
# {{ ansible_managed }}
# Please do not modify this file manually.

[Unit]
Description=Airflow webserver daemon
After=network.target

[Service]
Environment="AIRFLOW_HOME={{ airflow_install_directory }}"
User={{ airflow_user }}
Group={{ airflow_user }}
Type=simple
ExecStart=/usr/local/bin/airflow webserver --pid /run/airflow/webserver.pid
Restart=always
RestartSec=11s
PrivateTmp=true
RuntimeDirectory=airflow
RuntimeDirectoryMode=755

[Install]
WantedBy=multi-user.target
